# heartbeat/com.bareclaw.test-message.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bareclaw.test-message</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__TEST_SCRIPT__</string>
    </array>
    <key>StartInterval</key>
    <integer>30</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/bareclaw-test-message-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/bareclaw-test-message-stderr.log</string>
</dict>
</plist>
